# chef/cookbooks/persona-metrics/files/default/etc/cron.d/process_metrics
# m h dom mon dow user  command
30 4 * * * bid_metrics /opt/bid_metrics/bin/process_metrics.sh