# chef/cookbooks/persona-monitor/files/default/etc/cron.d/pagerduty_nagios
* * * * * nagios /usr/local/nagios/libexec/notifications/pagerduty_nagios.pl flush