# chef/cookbooks/bigtent/files/default/etc/init/svscan.conf
start on runlevel [345]
respawn
exec /command/svscanboot
